# deploy/safarai.service
[Unit]
Description=SafarAI Next.js Application
After=network.target postgresql.service

[Service]
Type=simple
User=safarai
WorkingDirectory=/opt/safarai
EnvironmentFile=/opt/safarai/.env.production
Environment=NODE_ENV=production
ExecStart=/usr/bin/node node_modules/next/dist/bin/next start -p 3000
Restart=on-failure
RestartSec=10

# Loglar journald ga
StandardOutput=journal
StandardError=journal
SyslogIdentifier=safarai

[Install]
WantedBy=multi-user.target
